# nmapp.service

[Unit]
Description=NetworkManagerApplet
After=wm.target

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=wm.target
