# notify.service

[Unit]
Description=Varslingstjeneste
After=wm.target

[Service]
ExecStart=/usr/lib/xfce4/notifyd/xfce4-notifyd

[Install]
WantedBy=wm.target
