# NumlockPaa.timer
[Unit]
Description=Kjører numlockx etter innlogg

[Timer]
OnActiveSec=10

[Install]
WantedBy=timers.target
