# NumlockPaa.service
[Unit]
Description=Numlock
After=wm.target

[Service]
ExecStart=/usr/bin/numlockx

[Install]
WantedBy=wm.target
