# udiskie.service

[Unit]
Description=Udiskie
After=wm.target

[Service]
ExecStart=/usr/bin/udiskie --tray 

[Install]
WantedBy=wm.target
